# systemd/mining-manager.service
[Unit]
Description=RTX 4080 Auto-Mining Manager (RainbowMiner)
After=network-online.target nvidia-persistenced.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/troggoman/mining-manager
ExecStart=/home/troggoman/mining-manager/mine.sh
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
Environment=DISPLAY=:0
Environment=XAUTHORITY=/root/.Xauthority

[Install]
WantedBy=multi-user.target
